# fan-controller.timer
[Unit]
Description=Run Fan Controller delayed 

[Timer]
OnBootSec=2min
Unit=fan-controller.service
[Install]
WantedBy=timers.target
